# templates/default/etc/systemd/system/etcd.service.erb
[Unit]
Description=etcd
Documentation=https://github.com/coreos
[Service]
Type=Notify
NotifyAccess=all
ExecStart=/usr/bin/etcd \
<%= render_str(@etcd_config['heartbeat-interval'],' --heartbeat-interval',' \\','=') -%>
<%= render_str(@etcd_config['cert-file'],' --cert-file',' \\','=') -%>
<%= render_str(@etcd_config['key-file'],' --key-file',' \\','=') -%>
<%= render_str(@etcd_config['peer-cert-file'],' --peer-cert-file',' \\','=') -%>
<%= render_str(@etcd_config['peer-key-file'],' --peer-key-file',' \\','=') -%>
<%= render_str(@etcd_config['trusted-ca-file'],' --trusted-ca-file',' \\','=') -%>
<%= render_str(@etcd_config['peer-trusted-ca-file'],' --peer-trusted-ca-file',' \\','=') -%>
<%= render_str(@etcd_config['initial-advertise-peer-urls'],' --initial-advertise-peer-urls',' \\','=') -%>
<%= render_str(@etcd_config['listen-peer-urls'],' --listen-peer-urls',' \\','=') -%>
<%= render_str(@etcd_config['listen-client-urls'],' --listen-client-urls',' \\','=') -%>
<%= render_str(@etcd_config['advertise-client-urls'],' --advertise-client-urls',' \\','=') -%>
<%= render_str(@etcd_config['initial-cluster-token'],' --initial-cluster-token',' \\','=') -%>
<%= render_str(@etcd_config['initial-cluster'],' --initial-cluster',' \\','=') -%>
<%= render_str(@etcd_config['initial-cluster-state'],' --initial-cluster-state',' \\','=') -%>
<%= render_str(@etcd_config['data-dir'],' --data-dir',' \\','=') -%>
<%= render_str(@etcd_config['name'],' --name','','=') -%>

Restart=on-failure
RestartSec=5
[Install]
WantedBy=multi-user.target 
