# templates/default/etc/systemd/system/kube-proxy.service.erb
[Unit]
Description=Kubernetes Kube Proxy
Documentation=https://github.com/GoogleCloudPlatform/kubernetes

[Service]
ExecStart=/usr/bin/kube-proxy \
<%= render_str(@kube_worker_config['kube-proxy']['kubeconfig'],' --kubeconfig',' \\','=') -%>
<%= render_str(@kube_worker_config['kube-proxy']['cluster-cidr'],' --cluster-cidr',' \\','=') -%>
<%= render_str(@kube_worker_config['kube-proxy']['proxy-mode'],' --proxy-mode',' \\','=') -%>
--v=2

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
